# rules/bootserver_dummywlan/templates/dummywlan.upstart.conf
description "Add dummy interface to wlan0 bridge"

start on static-network-up

task

exec /usr/local/lib/dummywlan
