# parts/wlan/debian/puavo-wlancontroller.upstart
description "puavo-wlancontroller"

start on runlevel [23]
stop on shutdown

chdir /opt/puavo-wlancontroller
env RACK_ENV=production
env PUAVO_WLANCONTROLLER_DB_SQLITE3=/var/opt/puavo-wlancontroller/db.sqlite3

exec sudo -E -u puavo-wlancontroller  bundle exec puma --workers 2 --port 45046

respawn
