# parts/wlan/debian/puavo-wlanap.upstart
description "Puavo's Tunneled WLAN Access Point"

start on (local-filesystems and net-device-up IFACE!=lo)

exec puavo-wlanap

respawn

post-stop exec sleep 60
