# parts/ltsp/bootserver/templates/etc/init/puavo-krb5-kdc.conf
description "puavo-krb5-kdc"

start on puavo-slapd-started
stop on puavo-slapd-stopped

expect fork

exec /usr/sbin/krb5kdc
