# parts/wlan/ap/etc/init.d/puavo-wlanap
#!/bin/sh

### BEGIN INIT INFO
# Provides:		puavo-wlanap
# Required-Start:	$remote_fs $syslog $network
# Required-Stop:	$remote_fs $syslog $network
# Should-Start:
# Should-Stop:
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Short-Description:	Puavo WLAN Access Point Daemon
# Description:		Puavo WLAN Access Point Daemon
### END INIT INFO

PATH=/sbin:/bin:/usr/sbin:/usr/bin
NAME=puavo-wlanap
EXECUTABLE=/usr/sbin/$NAME
DESC="Puavo WLAN Access Point Daemon"
PIDFILE=/var/run/$NAME.pid

[ -x "$EXECUTABLE" ] || exit 0

. /lib/lsb/init-functions

case "$1" in
  start)
	log_daemon_msg "Starting $DESC" "$NAME"
	start-stop-daemon --start --exec "$EXECUTABLE" \
            --pidfile "$PIDFILE" -- --daemon
	log_end_msg "$?"
	;;
  stop)
	log_daemon_msg "Stopping $DESC" "$NAME"
	start-stop-daemon --stop --pidfile "$PIDFILE"
	log_end_msg "$?"
	;;
  restart)
  	$0 stop
	sleep 3
	$0 start
	;;
  status)
	status_of_proc -p "$PIDFILE" "$EXECUTABLE" "$NAME"
	exit $?
	;;
  *)
	echo "Usage: service $NAME {start|stop|restart|status}" >&2
	exit 1
	;;
esac

exit 0
