# parts/sharedir/manager/systemd/puavo-sharedir-manager.service
[Unit]
Description=Puavo Sharedir Manager

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/sbin/puavo-sharedir-manager
