# rules/bootserver_network_interfaces/templates/etc/network/if-up.d/disable_flow_control
#!/bin/bash

set -eu

# XXX get rid of hosttype handling
puavo_hosttype=$(puavo-conf puavo.hosttype)

if [ "$puavo_hosttype" != 'bootserver' ]; then
  exit 0
fi

if [[ "$IFACE" =~ ^ltsp[0-9.]+$ ]]; then
  /sbin/ethtool -A $IFACE autoneg off rx off tx off
fi
