# rules/bootserver_network_interfaces/templates/etc/network/if-up.d/disable_flow_control
#!/bin/bash

if [[ "$IFACE" =~ ^eth[0-9.]+$ ]]; then
  /sbin/ethtool -A $IFACE autoneg off rx off tx off
fi
