# parts/ltsp/client/systemd/puavo-jetpipe-usb@.service
[Unit]
BindTo=%i.device
After=%i.device

[Service]
Type=oneshot
ExecStart=/usr/lib/puavo-ltsp-client/puavo-jetpipe usb /%I
