# cloud/puavo-ds/debian.default/puavo-ds-slave.puavo-krb5-kdc.upstart
description "puavo-krb5-kdc"

start on started slapd
stop on stopped slapd

respawn
exec /usr/sbin/krb5kdc -n
post-stop exec sleep 10
