# cloud/puavo-ds/debian.default/puavo-ds-slave.puavo-krb5-kdc.upstart
description "puavo-krb5-kdc"

start on puavo-slapd-started
stop on puavo-slapd-stopped

respawn
exec /usr/sbin/krb5kdc -n
post-stop exec sleep 10
