# parts/ltsp/client/templates/etc/network/if-up.d/dnssd-printers
#!/bin/sh

PUAVO_HOSTTYPE=$(cat /etc/puavo/hosttype)

if [ "${PUAVO_HOSTTYPE}" = "laptop" ]; then
  /usr/sbin/puavo-dnssd-printer-client >/dev/null 2>&1 &
fi

