# rules/bootserver_cups/files/cups-watchdog.service
[Unit]
Description=CUPS Watchdog
After=cups.service

[Service]
ExecStart=/usr/local/lib/cups-watchdog

[Install]
WantedBy=multi-user.target
