# parts/vpn/client/systemd/puavo-vpn-client-dnsmasq.service.in
[Unit]
Description=Puavo VPN Client (Dnsmasq)

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=#FILLED BY MAKEFILE DURING INSTALL
