# parts/wlan/ap/puavo-wlanap.service.in
[Unit]
Description=Puavo WLAN Access Point
After=network.target
Wants=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=#FILLED BY MAKEFILE DURING INSTALL
