# rules/docker/files/puavo-docker.timer
[Unit]
Description=Manage Docker containers

[Timer]
OnStartupSec=5m
OnCalendar=02:47

[Install]
WantedBy=timers.target
