# rules/bootserver_autopoweron/templates/puavo-autopoweron.cron
30 07 * * * root /sbin/start --quiet puavo-autopoweron >/dev/null 2>&1
