# rules/bootserver_cron/files/fix-homedir-permissions.cron
# Fix home directory permissions every morning at 06:13
13 06 * * * root /usr/local/sbin/puavo-fix-homedir-permissions
