# rules/bootserver_cron/files/fstrim.cron
PATH=/bin:/sbin

# Trim SSD-cached /home every night at 03:00
00 03 * * * root grep -q '^/dev/mapper/cache-home /home ext4' && fstrim /home >/dev/null
