# parts/ltsp/client/templates/etc/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/lib/puavo-ltsp-client/hotplug-monitor"
