# parts/ltsp/client/etc/network/if-up.d/sssd
#!/bin/sh

# Don't do anything for the lo interface
test "$IFACE" != lo || exit 0

PUAVO_HOSTTYPE=$(cat /etc/puavo/hosttype)

if [ "$PUAVO_HOSTTYPE" != "laptop" ]; then
  exit 0
fi

# Make sure that sssd goes online
sleep 1
killall -USR2 sssd
sleep 6
killall -USR2 sssd
