# debs/autopoweroff/debian/autopoweroff.service
[Unit]
Description=Autopoweroff

[Service]
ExecStart=/usr/sbin/autopoweroffd

[Install]
WantedBy=multi-user.target
