# parts/vpn/client/dispatcher.d/puavo-vpn-dnsmasq
#!/bin/sh

/usr/lib/puavo-vpn-client/set-servers
