# parts/ltsp/puavo-install/upstart/puavo-update-client.conf
# puavo-update-client - Update Puavo client (locally installed)
#
# This task updates this host a new operating system version, plus handles some
# configuration updates.

description	"Update Puavo host (locally installed)"

start on net-device-up IFACE!=lo
stop on runlevel [016]

script
  random_number=$(printf %d "0x$(openssl rand -hex 2)")

  # Wait some time between one and five minutes.
  wait_seconds=$(expr "$random_number" % 240 + 60)
  sleep $wait_seconds

  # The real work is done elsewhere.  Trigger puavo-update-client and wait
  # for an hour before triggering it again.
  while true; do
    puavo-update-client --use-rate-limit
    sleep 3600
  done
end script
