# parts/ltsp/puavo-install/upstart/puavo-update-laptop.conf
# puavo-update-laptop - Update Puavo laptop
#
# This task updates this host a new operating system version, plus handles some
# configuration updates.

description	"Update Puavo laptop"

start on net-device-up
stop on runlevel [016]

task

script
  puavo-handle-external-files-actions

  random_number=$(printf %d "0x$(openssl rand -hex 2)")

  # wait some time between one and five minutes
  wait_seconds=$(expr "$random_number" % 240 + 60)
  sleep $wait_seconds

  # the real work is done elsewhere
  puavo-update-laptop
end script
