# parts/ltsp/puavo-install/upstart/puavo-update-client.conf
# puavo-update-client - Update Puavo client (locally installed)
#
# This task updates this host a new operating system version, plus handles some
# configuration updates.

description	"Update Puavo host (locally installed)"

start on net-device-up IFACE!=lo and started dbus
stop on runlevel [016] or stopping dbus

script
  # if booted from nbd, do nothing
  if grep -q root=/dev/nbd /proc/cmdline; then
    stop; exit 0
  fi

  # Wait 90 seconds before proceeding (seems like a nice number, at least
  # the administrative tunnel should be up by then).
  sleep 90

  # The real work is done elsewhere.  Trigger puavo-update-client and wait
  # for an hour before triggering it again.
  while true; do
    # If automatic updates are false or we are in development mode (persistent
    # overlay is in use), do not do an image update.
    automatic_image_updates=$(jq -r .automatic_image_updates \
                                    /etc/puavo/device.json)
    if [ "${automatic_image_updates}" = "false" ] \
         || grep -q "puavo.image.overlay=" /proc/cmdline; then
      puavo-update-client --no-image-update || true
    else
      # Do image update automatically, but do also use rate-limit.
      puavo-update-client --use-rate-limit || true
    fi
    sleep 3600
  done
end script
