# rules/bootserver_ddns/files/puavo-update-airprint-ddns.cron
*/15 * * * * root /usr/local/sbin/puavo-update-airprint-ddns >/dev/null 2>&1
