# parts/ltsp/puavo-install/systemd/puavo-torrent-updated.socket
[Socket]
ListenStream=/run/puavo/puavo-torrent-updated.sock
SocketMode=0600
