# parts/ltsp/client/gdm/greeter/autostart/puavo-run-xrandr.desktop
[Desktop Entry]
Type=Application
Name=Puavo run xrandr
Exec=/usr/lib/puavo-ltsp-client/run-xrandr
NoDisplay=true
