# parts/ltsp/debian.default/puavo-ltsp-bootserver.puavo-update-airprint-ddns.cron.d
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

*/15 * * * * root puavo-update-airprint-ddns >/dev/null 2>&1
