# rules/homedir_management/files/puavo-cleanup-old-users.cron
#!/bin/sh

test -x /usr/local/sbin/puavo-cleanup-old-users || exit 0
/usr/local/sbin/puavo-cleanup-old-users
