# parts/hw-log/debian.default/puavo-hw-log.upstart
description "puavo-hw-log"

start on runlevel [23]

exec sh -c 'sleep 10 && puavo-hw-log'

