# parts/hw-log/debian.default/puavo-hw-log.upstart
description "puavo-hw-log"

start on net-device-up IFACE!=lo
stop on runlevel [016]

script
  sleep 10
  puavo-hw-log
end script
