# parts/ltsp/client/upstart/puavo-report-system-version.conf
description "Puavo report system version to syslog"

start on runlevel [2345]

task
exec puavo-report-system-version
