# parts/ltsp/client/upstart/puavo-sync-external-files.conf
description "Sync external files from Puavo"

start on runlevel [2345]

task

# for laptops this is done through puavo-update-laptop
exec puavo-sync-external-files
