# parts/ltsp/bootserver/upstart/puavo-redirect-freenx.conf
description "puavo-redirect-freenx"

start on runlevel [2345]
stop on runlevel [016]

exec /usr/sbin/puavo-redirect-freenx
