# parts/autopoweroff/puavo-autopoweroff.service
[Unit]
Description=Puavo autopoweroff

[Service]
ExecStart=/usr/sbin/puavo-autopoweroff

[Install]
WantedBy=multi-user.target
