# parts/aptirepo/debian.default/debbox.upstart
description "debbox"

start on runlevel [23]
stop on shutdown

env NODE_ENV=production
exec sudo -E -u www-data debbox

respawn
