# parts/ltsp/client/etc/NetworkManager/dispatcher.d/dnssd-printers
#!/bin/sh

IFACE=$1
COMMAND=$2

if [ "${COMMAND}" = "up" ]; then
  PUAVO_HOSTTYPE=$(cat /etc/puavo/hosttype)

  if [ "${PUAVO_HOSTTYPE}" = "laptop" ]; then
    /usr/sbin/puavo-dnssd-printer-client --retries 6 --retry-delay 15 --delay 10 >/dev/null 2>&1 &
  fi
fi
