# rules/wacom/files/60-wacom.rules
SUBSYSTEM=="usb", ACTION=="add", ATTRS{idVendor}=="056a", ATTRS{idProduct}=="0303", RUN+="/sbin/modprobe wacom"
