# rules/abitti/files/puavo-update-abitti-image.cron
#!/bin/sh

test -x /usr/local/sbin/puavo-update-abitti-image || exit 0

exec /usr/local/sbin/puavo-update-abitti-image --use-rate-limit
