# parts/rules/ltsp-build-client/puppet/ltsp/apt_ltspbuild_proxy/templates/remove-apt-proxy.conf
# remove-apt-proxy

description     "remove apt proxy setting at boot"

start on filesystem

task
script
  /bin/rm -f /etc/apt/apt.conf.d/00ltspbuild-proxy
end script
