# parts/rules/ltsp-build-client/puppet/ltsp/polyvision/templates/eno.upstart
description "eno"
 
start on runlevel [23]
stop on shutdown
 
script
    exec /sbin/eno-connect
end script
 
respawn
