# parts/displays/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/lib/puavo-displays/udev-hotplug-monitor"
