# rules/bootserver_cups/templates/cups-watchdog.upstart
description "cups-watchdog"

start on started cups
stop on stopped cups

exec /usr/local/lib/cups-watchdog

respawn
