# parts/sharedir/upstart/puavo_share_manager.conf
# puavo_share_manager - Manage Puavo Shares

description     "PuavoShare Manager"

start on filesystem
stop on runlevel S

respawn
respawn limit 10 5

exec /usr/local/sbin/puavo_share_manager
